# scripts/trading-evaluator-task.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Jarvis Trading Bot Evaluator — taeglich nach US Market Close</Description>
    <Author>Robin Mandel</Author>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2026-04-21T22:05:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <ExecutionTimeLimit>PT10M</ExecutionTimeLimit>
    <Enabled>true</Enabled>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>C:\Users\Robin\AppData\Local\Programs\Python\Python311\python.exe</Command>
      <Arguments>C:\Users\Robin\Jarvis\trading-bots\evaluator.py</Arguments>
      <WorkingDirectory>C:\Users\Robin\Jarvis\trading-bots</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
